# services/controller.service
[Unit]
Description=WATCHES Fan Controller

[Service]
ExecStart=
WorkingDirectory=
Restart=always
User=

[Install]
WantedBy=multi-user.target