# services/sensor.service
[Unit]
Description=WATCHES Sensor Driver

[Service]
ExecStart=
WorkingDirectory=
Restart=always
User=

[Install]
WantedBy=multi-user.target