# services/server.service
[Unit]
Description=WATCHES Server

[Service]
ExecStart=
WorkingDirectory=
Restart=always
User=

[Install]
WantedBy=multi-user.target